# ansible/ELK6/playbooks/roles/opendistro/files/opendistro.service
[Unit]
Description=opendistro

[Service]
ExecStart=/usr/local/bin/docker-compose -f /usr/local/etc/opendistro/docker-compose.yml up
ExecStop=/usr/local/bin/docker-compose -f /usr/local/etc/opendistro/docker-compose.yml down
ExecReload=/usr/local/bin/docker-compose -f /usr/local/etc/opendistro/docker-compose.yml restart
Restart=always

Type=simple

[Install]
WantedBy=multi-user.target
